# spark/files/systemd/hive-server2.service.d/default.conf
[Service]
Environment="JAVA_HOME=/usr/lib/jvm/java"
Environment="PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin"
