# spark/files/systemd/phoenix-query-server.service.d/default.conf
[Service]
Environment="JAVA_HOME=/usr/lib/jvm/java"
Environment="PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin"
Environment="HBASE_CONF_PATH=/etc/hbase/conf"
