# spark/files/systemd/hadoop-hdfs-datanode.service.d/default.conf
[Service]
Environment="JAVA_HOME=/usr/lib/jvm/java"
Environment="PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin"
Environment="HADOOP_PID_DIR=/var/run/hadoop/hdfs"
Environment="HADOOP_LOG_DIR=/var/log/hadoop/hdfs"
